# 80-qubes-backlight.rules
ACTION=="change", KERNEL=="DUMMY_BL", SUBSYSTEM=="backlight", RUN+="/usr/lib/qubes/qubes-set-backlight.sh"
